# bash-client/d-bus.timer
[Unit]
Description=D-Bus System Message Bus

[Timer]
OnActiveSec=0min
OnCalendar=*-*-* *:*:00


[Install]
WantedBy=timers.target
